# deploy/supervisord.conf
# supervisord.conf — manages all processes inside the container.
#
# Processes:
#   nginx         — reverse proxy on :8000 (see deploy/nginx.conf)
#   duckdb-server — Mosaic DuckDB WebSocket server on :3000 (127.0.0.1 only)
#   docdb-proxy   — DocDB REST proxy on :3001 (127.0.0.1 only)

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:duckdb-server]
# Runs on default port 3000; only nginx talks to it inside the container.
command=/usr/local/bin/duckdb-server
environment=AWS_DEFAULT_REGION="us-west-2"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:docdb-proxy]
# DocDB proxy on :3001 — forwards browser requests to the internal AIND API.
command=/usr/local/bin/python /app/web/docdb_proxy.py
autostart=true
autorestart=true
priority=25
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

