# deploy/supervisord.conf
# supervisord.conf — manages all three processes inside the container.
#
# Processes:
#   nginx         — reverse proxy on :8000 (see deploy/nginx.conf)
#   duckdb-server — Mosaic DuckDB WebSocket server on :3000 (127.0.0.1 only)
#   panel         — Panel/Bokeh server on :8001 (127.0.0.1 only)
#
# Environment variables consumed at runtime:
#   ALLOW_WEBSOCKET_ORIGIN   passed through to `panel serve --allow-websocket-origin`
#   OAUTH_REDIRECT           passed through to `panel serve --oauth-redirect-uri`

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:duckdb-server]
# Runs on default port 3000; only nginx talks to it inside the container.
command=/usr/local/bin/duckdb-server
environment=AWS_DEFAULT_REGION="us-west-2"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:panel]
command=/bin/sh -c "panel serve /app/src/zombie/assets.py /app/src/zombie/contributions.py /app/src/zombie/subject.py --static-dirs images=/app/src/zombie/images --address 127.0.0.1 --port 8001 --allow-websocket-origin $ALLOW_WEBSOCKET_ORIGIN --oauth-redirect-uri $OAUTH_REDIRECT --keep-alive 10000 --num-threads 4"
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
